# node-confs/jumplogs.conf
description "jumphost logs collector"
start on runlevel [23]
start on started tty1

respawn
respawn limit 10 5
#expect fork
exec sudo -u jumplogs /usr/local/bin/asciinema-process.sh 2>&1
