# node-confs/restapi.conf
description "kvm rest api"
start on runlevel [23]
start on started tty1

respawn
respawn limit 10 5
#expect fork
exec sudo -u restapi /home/restapi/runserver.sh 2>&1 | logger -t restapi

